# deploy/systemd/zaehler.service
# systemd-Unit für die Zählerstand-App.
#
# Pfade nehmen an, dass die Anwendung unter /opt/zaehler installiert ist und
# der App-Benutzer "zaehler" heißt. Anpassen, falls anders deployed.
# Aktivierung: siehe deploy/lxc/README.md, Abschnitt 4.6.

[Unit]
Description=Zaehlerstand-App (FastAPI + Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zaehler
Group=zaehler
WorkingDirectory=/opt/zaehler/backend

# uv und pnpm liegen im User-lokalen .local/bin
Environment=PATH=/opt/zaehler/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

# Standardpfad für die SQLite-Datenbank. Der Pfad wird durch meters.env
# überschrieben werden, falls die Datei existiert.
Environment=METERS_DATABASE_URL=sqlite:////opt/zaehler/data/meters.db

# Optionale Konfigurationsdatei (METERS_SECRET_KEY, METERS_COOKIE_SECURE, ...).
# Das führende '-' weist systemd an, das Fehlen NICHT als Fehler zu werten.
EnvironmentFile=-/opt/zaehler/data/meters.env

# Bind-Host und -Port kommen aus meters.env (METERS_BIND_HOST, METERS_BIND_PORT).
# Default: 0.0.0.0:8000 — direkter LAN-Zugriff.
# Sobald ein HTTPS-Reverse-Proxy davor steht, in meters.env auf 127.0.0.1 setzen.
ExecStart=/opt/zaehler/.local/bin/uv run python -m meters.server

Restart=on-failure
RestartSec=3

# --- Hardening (siehe systemd.exec(5)) ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
# Schreibzugriff auf das gesamte App-Verzeichnis — uv legt dort seinen Cache
# (~/.cache/uv) und die virtuelle Umgebung (backend/.venv) an, der Frontend-
# Build schreibt nach backend/src/meters/static/. Außerhalb von /opt/zaehler
# bleibt das Dateisystem read-only.
ReadWritePaths=/opt/zaehler
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
CapabilityBoundingSet=
AmbientCapabilities=
PrivateUsers=true
UMask=0027

[Install]
WantedBy=multi-user.target
